# etc/cron.d/cronfile.stub
#!/usr/bin/env bash
# $Id: monthly,v 1.6 2004/10/03 08:20:01 davep Exp $
#set -x

# set up more featureful environment than CRON does
# XXX do this in CRONTAB

. $HOME/etc/crontab.rc

true
